# errorReporting/crontab
@reboot /home/ubuntu/IT490/errorReporting/logginserver.php